# packages/cos-setup/cos-setup-network.openrc
#!/sbin/openrc-run

depend() {
  after net
  provide cos-setup-network
}

start() {
  cos-setup network
  eend 0
}